# templates/default/consul.conf.erb
description     "Consul agent"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
  OPTIONS=""
  if [ -e /etc/default/consul ]
  then
    . /etc/default/consul
    echo "Sourcing config file."
  else
    echo "Config file not available."
  fi
  exec /usr/local/bin/consul agent -config-file <%= node['consul']['config_file_path'] %> $OPTIONS
end script
